# .config/systemd/user/junest.service
[Unit]
Description=Create JuNest image

[Service]
ExecStart=/home/builder/build_image.sh
Type=simple

[Install]
WantedBy=default.target

